# in/delete-old-files.cron
0 * * * * root /usr/bin/delete-old-files.sh